# runtime-files/etc/services.d/10-syslog/run
#!/bin/bash

#
# Starts syslog
# =============
#

set -e

. /usr/local/lib/env.inc

# Load the runtime environment
load_env

exec s6-envdir ${ENV_DIR} syslogd -O - -t -R $HOST_IP:514 -L -n
